# output-baseline/VirusTotal-behaviors/eae3edbba73566794b10e2c485b520e429420b529aad8bdaeb2309ec4bc5e522.json
{
    "data":  {
                 "attributes":  {
                                    "registry_keys_set":  [
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\GoogleUpdaterInternalService126.0.6441.0\\Start = DWORD (0x00000002)",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\GoogleUpdaterInternalService126.0.6441.0\\ImagePath = \"C:\\Program Files (x86)\\Google\\GoogleUpdater\\126.0.6441.0\\updater.exe\" --system --windows-service --service=update-internal",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\GoogleUpdaterInternalService126.0.6462.0\\Start = DWORD (0x00000002)",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\GoogleUpdaterInternalService126.0.6462.0\\ImagePath = \"C:\\Program Files (x86)\\Google\\GoogleUpdater\\126.0.6462.0\\updater.exe\" --system --windows-service --service=update-internal",
                                                              "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Policies\\Microsoft\\Windows\\Safer\\CodeIdentifiers",
                                                              "\\Registry\\Machine\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\996E.exe",
                                                              "\\Registry\\MACHINE\\System\\CurrentControlSet\\Control\\SafeBoot\\Option",
                                                              "\\REGISTRY\\MACHINE\\SOFTWARE\\Policies\\Microsoft\\Windows\\Safer\\CodeIdentifiers\\TransparentEnabled",
                                                              "\\Registry\\Machine\\Software\\Policies\\Microsoft\\Windows\\Safer\\CodeIdentifiers",
                                                              "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Explorer\\Shell Folders",
                                                              "HKEY_LOCAL_MACHINE\\SYSTEM\\ControlSet001\\Control\\MUI\\UILanguages\\en-US",
                                                              "HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Control\\Nls\\CustomLocale",
                                                              "HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Control\\Session Manager",
                                                              "HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Control\\Session Manager\\Segment Heap",
                                                              "HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows NT\\CurrentVersion\\AppCompatFlags\\Disable8And16BitMitigation",
                                                              "HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Wow64\\x86",
                                                              "HKEY_LOCAL_MACHINE\\Software\\WOW6432Node\\Microsoft\\Windows NT\\CurrentVersion\\AppCompatFlags",
                                                              "HKEY_LOCAL_MACHINE\\Software\\WOW6432Node\\Policies\\Microsoft\\MUI\\Settings",
                                                              "HKEY_LOCAL_MACHINE\\Software\\WOW6432Node\\Policies\\Microsoft\\Windows\\Safer\\CodeIdentifiers",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\FileSystem\\",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\MUI\\Settings\\LanguageConfiguration",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\MUI\\UILanguages",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\MUI\\UILanguages\\PendingDelete",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\NLS\\Language",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\SafeBoot\\Option",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\Session Manager",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\Srp\\GP\\DLL",
                                                              "HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Control\\Terminal Server"
                                                          ],
                                    "mutexes_created":  [

                                                        ],
                                    "processes_created":  [
                                                              "%SAMPLEPATH%\\lua.exe",
                                                              "C:\\Windows\\System32\\wuapihost.exe",
                                                              "C:\\Program Files (x86)\\Google3656_2088715030\\bin\\updater.exe",
                                                              "%SAMPLEPATH%\\eae3edbba73566794b10e2c485b520e429420b529aad8bdaeb2309ec4bc5e522.exe",
                                                              "C:\\Windows\\System32\\UI0Detect.exe",
                                                              "C:\\Program Files\\Google3980_1061348839\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google3264_76114242\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google2500_1976639986\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google2740_550104856\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google2744_2022486992\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google1276_757632363\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google1868_1417054015\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google2424_1175041261\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google3972_1953574762\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google3860_620853748\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google2484_932641056\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google1556_1478716743\\bin\\updater.exe",
                                                              "C:\\Program Files\\Google2076_1401794204\\bin\\updater.exe"
                                                          ],
                                    "registry_keys_opened":  [

                                                             ],
                                    "mitre_attack_techniques":  [
                                                                    {
                                                                        "id":  "T1129",
                                                                        "severity":  "IMPACT_SEVERITY_INFO",
                                                                        "signature_description":  "link function at runtime on Windows",
                                                                        "refs":  [
                                                                                     {
                                                                                         "ref":  "#signature_matches",
                                                                                         "value":  "5"
                                                                                     }
                                                                                 ]
                                                                    },
                                                                    {
                                                                        "id":  "T1129",
                                                                        "severity":  "IMPACT_SEVERITY_INFO",
                                                                        "signature_description":  "link many functions at runtime",
                                                                        "refs":  [
                                                                                     {
                                                                                         "ref":  "#signature_matches",
                                                                                         "value":  "9"
                                                                                     }
                                                                                 ]
                                                                    },
                                                                    {
                                                                        "id":  "T1071",
                                                                        "signature_description":  "Adversaries may communicate using application layer protocols to avoid detection/network filtering by blending in with existing traffic."
                                                                    },
                                                                    {
                                                                        "signature_description":  "Tries to load missing DLLs",
                                                                        "id":  "T1574.002",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Spawns processes",
                                                                        "id":  "T1055",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Reads software policies",
                                                                        "id":  "T1082",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "May sleep (evasive loops) to hinder dynamic analysis",
                                                                        "id":  "T1497",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Sample monitors Window changes (e.g. starting applications), analyze the sample with the simulation cookbook",
                                                                        "id":  "T1010",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Executes the \"sed\" command used to modify input streams (typically from files or pipes)",
                                                                        "id":  "T1059",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Executes commands using a shell command-line interpreter",
                                                                        "id":  "T1064",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Executes the \"systemctl\" command used for controlling the systemd system and service manager",
                                                                        "id":  "T1543.002",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Creates hidden files, links and/or directories",
                                                                        "id":  "T1564.001",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Uses the \"uname\" system call to query kernel version information (possible evasion)",
                                                                        "id":  "T1518.001",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Sample reads /proc/mounts (often used for finding a writable filesystem)",
                                                                        "id":  "T1083",
                                                                        "severity":  "IMPACT_SEVERITY_LOW"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Reads system information from the proc file system",
                                                                        "id":  "T1082",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "signature_description":  "Reads CPU information from /sys indicative of miner or evasive malware",
                                                                        "id":  "T1082",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "refs":  [
                                                                                     {
                                                                                         "ref":  "#signature_matches",
                                                                                         "value":  "187"
                                                                                     }
                                                                                 ],
                                                                        "signature_description":  "Tries to load missing DLLs",
                                                                        "id":  "T1574.002",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "refs":  [
                                                                                     {
                                                                                         "ref":  "#signature_matches",
                                                                                         "value":  "311"
                                                                                     }
                                                                                 ],
                                                                        "signature_description":  "Creates a DirectInput object (often for capturing keystrokes)",
                                                                        "id":  "T1056",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    },
                                                                    {
                                                                        "refs":  [
                                                                                     {
                                                                                         "ref":  "#signature_matches",
                                                                                         "value":  "509"
                                                                                     }
                                                                                 ],
                                                                        "signature_description":  "Reads software policies",
                                                                        "id":  "T1082",
                                                                        "severity":  "IMPACT_SEVERITY_INFO"
                                                                    }
                                                                ]
                                },
                 "id":  "eae3edbba73566794b10e2c485b520e429420b529aad8bdaeb2309ec4bc5e522",
                 "type":  "file_behaviour_summary"
             }
}
