# output-baseline/VirusTotal-behaviors/malicious/047bcdfda2ed0aea64d35fa608962fb3f7dc21eda36bc75b083a30bae3d1fcc5.json
{
    "data":  {
                 "verdict_labels":  [
                                        "KeyloggerGeneric"
                                    ],
                 "mbc":  [
                             {
                                 "id":  "B0013.001",
                                 "objective":  "Discovery",
                                 "behavior":  "Analysis Tool Discovery",
                                 "method":  "Process detection",
                                 "refs":  [
                                              {
                                                  "ref":  "#signature_matches",
                                                  "value":  "1"
                                              }
                                          ]
                             },
                             {
                                 "id":  "F0012",
                                 "objective":  "Persistence",
                                 "behavior":  "Registry Run Keys / Startup Folder",
                                 "refs":  [
                                              {
                                                  "ref":  "#signature_matches",
                                                  "value":  "2"
                                              }
                                          ]
                             },
                             {
                                 "id":  "B0023",
                                 "objective":  "Execution",
                                 "behavior":  "Install Additional Program",
                                 "refs":  [
                                              {
                                                  "ref":  "#signature_matches",
                                                  "value":  "3"
                                              }
                                          ]
                             }
                         ],
                 "processes_tree":  [
                                        {
                                            "process_id":  "2744",
                                            "name":  "%windir%\\explorer.exe %ProgramFiles(x86)%\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\SVE\\AdobeID.exe"
                                        },
                                        {
                                            "process_id":  "2232",
                                            "name":  "%windir%\\System32\\svchost.exe -k WerSvcGroup"
                                        },
                                        {
                                            "process_id":  "2600",
                                            "name":  "%TEMP%\\CD9J6RO4AHMVRLGT.exe"
                                        }
                                    ],
                 "signature_matches":  [
                                           {
                                               "id":  "1",
                                               "format":  "SIG_FORMAT_CAPA",
                                               "name":  "reference analysis tools strings",
                                               "description":  "anti-analysis",
                                               "authors":  [
                                                               "michael.hunhoff@mandiant.com"
                                                           ],
                                               "match_data":  [
                                                                  "regex: /(?\u003c!\\w)ida[gqtuw]?(\\.exe)?$/i",
                                                                  "regex: /windbg(\\.exe)?/i"
                                                              ],
                                               "rule_src":  "rule:\n  meta:\n    name: reference analysis tools strings\n    namespace: anti-analysis\n    authors:\n      - michael.hunhoff@mandiant.com\n    scopes:\n      static: file\n      dynamic: file\n    mbc:\n      - Discovery::Analysis Tool Discovery::Process detection [B0013.001]\n    references:\n      - https://github.com/LordNoteworthy/al-khaser/blob/master/al-khaser/AntiAnalysis/process.cpp\n    examples:\n      - al-khaser_x86.exe_\n  features:\n    - or:\n      - string: /ollydbg(\\.exe)?/i\n      - string: /ProcessHacker(\\.exe)?/i\n      - string: /tcpview(\\.exe)?/i\n      - string: /autoruns(\\.exe)?/i\n      - string: /autorunsc(\\.exe)?/i\n      - string: /filemon(\\.exe)?/i\n      - string: /procmon(\\.exe)?/i\n      - string: /regmon(\\.exe)?/i\n      - string: /procexp(\\.exe)?/i\n      - string: /(?\u003c!\\w)ida[gqtuw]?(\\.exe)?$/i\n      - string: /ida[gqtuw]?64(\\.exe)?$/i\n      - string: /ImmunityDebugger(\\.exe)?/i\n      - string: /Wireshark(\\.exe)?/i\n      - string: /dumpcap(\\.exe)?/i\n      - string: /HookExplorer(\\.exe)?/i\n      - string: /ImportREC(\\.exe)?/i\n      - string: /PETools(\\.exe)?/i\n      - string: /LordPE(\\.exe)?/i\n      - string: /SysInspector(\\.exe)?/i\n      - string: /proc_analyzer(\\.exe)?/i\n      - string: /sysAnalyzer(\\.exe)?/i\n      - string: /sniff_hit(\\.exe)?/i\n      - string: /windbg(\\.exe)?/i\n      - string: /joeboxcontrol(\\.exe)?/i\n      - string: /joeboxserver(\\.exe)?/i\n      - string: /ResourceHacker(\\.exe)?/i\n      - string: /x32dbg(\\.exe)?/i\n      - string: /x64dbg(\\.exe)?/i\n      - string: /Fiddler(\\.exe)?/i\n      - string: /httpdebugger(\\.exe)?/i\n      - string: /fakenet(\\.exe)?/i\n      - string: /netmon(\\.exe)?/i\n      - string: /WPE PRO(\\.exe)?/i\n      - string: /decompile(\\.exe)?/i\n      - string: /scylla/i\n      - string: /megadumper/i\n",
                                               "refs":  [
                                                            {
                                                                "ref":  "https://github.com/LordNoteworthy/al-khaser/blob/master/al-khaser/AntiAnalysis/process.cpp"
                                                            }
                                                        ]
                                           },
                                           {
                                               "id":  "2",
                                               "format":  "SIG_FORMAT_CAPA",
                                               "name":  "persist via Run registry key",
                                               "description":  "persistence/registry/run",
                                               "authors":  [
                                                               "moritz.raabe@mandiant.com",
                                                               "mehunhoff@google.com"
                                                           ],
                                               "match_data":  [
                                                                  "regex: /Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run/i",
                                                                  "HKEY_CURRENT_USER"
                                                              ],
                                               "rule_src":  "rule:\n  meta:\n    name: persist via Run registry key\n    namespace: persistence/registry/run\n    authors:\n      - moritz.raabe@mandiant.com\n      - mehunhoff@google.com\n    scopes:\n      static: function\n      dynamic: span of calls\n    att\u0026ck:\n      - Persistence::Boot or Logon Autostart Execution::Registry Run Keys / Startup Folder [T1547.001]\n    mbc:\n      - Persistence::Registry Run Keys / Startup Folder [F0012]\n    examples:\n      - Practical Malware Analysis Lab 06-03.exe_:0x401130\n      - b87e9dd18a5533a09d3e48a7a1efbcf6:0x1400070E0\n      - 9ff8e68343cc29c1036650fc153e69f7:0x470624\n  features:\n    - and:\n      - or:\n        - match: set registry value\n        - number: 0x80000001 = HKEY_CURRENT_USER\n        - number: 0x80000002 = HKEY_LOCAL_MACHINE\n      - or:\n        - string: /Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run/i\n        - string: /Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\Shell Folders/i\n        - string: /Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\User Shell Folders/i\n        - string: /Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunServices/i\n        - string: /Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Policies\\\\Explorer\\\\Run/i\n        - string: /Software\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Windows\\\\Load/i\n        - string: /System\\\\(ControlSet\\d{3}|CurrentControlSet)\\\\Control\\\\Session Manager\\\\BootExecute/i\n"
                                           },
                                           {
                                               "id":  "3",
                                               "format":  "SIG_FORMAT_CAPA",
                                               "name":  "contain an embedded PE file",
                                               "description":  "executable/subfile/pe",
                                               "authors":  [
                                                               "moritz.raabe@mandiant.com"
                                                           ],
                                               "rule_src":  "rule:\n  meta:\n    name: contain an embedded PE file\n    namespace: executable/subfile/pe\n    authors:\n      - moritz.raabe@mandiant.com\n    scopes:\n      static: file\n      dynamic: file\n    mbc:\n      - Execution::Install Additional Program [B0023]\n    examples:\n      - Practical Malware Analysis Lab 01-04.exe_:0x4060\n  features:\n    - or:\n      - count(characteristic(embedded pe)): 1 or more\n      - count(string(This program cannot be run in DOS mode.)): 2 or more\n"
                                           },
                                           {
                                               "id":  "5",
                                               "format":  "SIG_FORMAT_CAPA",
                                               "name":  "compiled from Visual Basic",
                                               "description":  "compiler/vb",
                                               "authors":  [
                                                               "@williballenthin"
                                                           ],
                                               "match_data":  [
                                                                  "regex: /VB5!.*/"
                                                              ],
                                               "rule_src":  "rule:\n  meta:\n    name: compiled from Visual Basic\n    namespace: compiler/vb\n    authors:\n      - \"@williballenthin\"\n    scopes:\n      static: file\n      dynamic: file\n    examples:\n      - 9bca6b99e7981208af4c7925b96fb9cf\n  features:\n    - and:\n      - string: /VB5!.*/\n      - import: msvbvm60.ThunRTMain\n"
                                           },
                                           {
                                               "id":  "6",
                                               "format":  "SIG_FORMAT_CAPA",
                                               "name":  "contains PDB path",
                                               "description":  "executable/pe/pdb",
                                               "authors":  [
                                                               "moritz.raabe@mandiant.com"
                                                           ],
                                               "match_data":  [
                                                                  "regex: /:\\\\.*\\.pdb/"
                                                              ],
                                               "rule_src":  "rule:\n  meta:\n    name: contains PDB path\n    namespace: executable/pe/pdb\n    authors:\n      - moritz.raabe@mandiant.com\n    scopes:\n      static: file\n      dynamic: file\n    examples:\n      - 464EF2CA59782CE697BC329713698CCC  # level32.exe\n  features:\n    - string: /:\\\\.*\\.pdb/\n"
                                           }
                                       ],
                 "verdict_confidence":  50,
                 "verdicts":  [
                                  "STEALER"
                              ],
                 "mitre_attack_techniques":  [
                                                 {
                                                     "id":  "T1547.001",
                                                     "severity":  "IMPACT_SEVERITY_INFO",
                                                     "signature_description":  "persist via Run registry key",
                                                     "refs":  [
                                                                  {
                                                                      "ref":  "#signature_matches",
                                                                      "value":  "2"
                                                                  }
                                                              ]
                                                 }
                                             ],
                 "attack_techniques":  {
                                           "T1547.001":  [
                                                             {
                                                                 "severity":  "INFO",
                                                                 "description":  "persist via Run registry key",
                                                                 "refs":  [
                                                                              {
                                                                                  "ref":  "#signature_matches",
                                                                                  "value":  "2"
                                                                              }
                                                                          ]
                                                             }
                                                         ]
                                       }
             }
}
